# docker/site-template/conf/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:php-fpm]
command=php-fpm --nodaemonize
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/php-fpm.err
stdout_logfile=/var/log/supervisor/php-fpm.out
priority=10

[program:nginx]
command=nginx -g "daemon off;"
autostart=%(ENV_WEB_SERVER_IS_NGINX)s
autorestart=true
stderr_logfile=/var/log/supervisor/nginx.err
stdout_logfile=/var/log/supervisor/nginx.out
priority=20

[program:apache2]
command=apache2ctl -DFOREGROUND
autostart=%(ENV_WEB_SERVER_IS_APACHE)s
autorestart=true
stderr_logfile=/var/log/supervisor/apache2.err
stdout_logfile=/var/log/supervisor/apache2.out
priority=20

[program:lighttpd]
command=lighttpd -D -f /etc/lighttpd/lighttpd.conf
autostart=%(ENV_WEB_SERVER_IS_LIGHTTPD)s
autorestart=true
stderr_logfile=/var/log/supervisor/lighttpd.err
stdout_logfile=/var/log/supervisor/lighttpd.out
priority=20

[program:node]
command=pm2-runtime start %(ENV_NODE_APP)s --name site-app
autostart=%(ENV_NODE_ENABLED)s
autorestart=true
stderr_logfile=/var/log/supervisor/node.err
stdout_logfile=/var/log/supervisor/node.out
priority=30

[program:sshd]
command=/usr/sbin/sshd -D -e
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/sshd.err
stdout_logfile=/var/log/supervisor/sshd.out
priority=5

[program:container-api]
; Runs as gnukontrolr-admin — manages the secure area (/var/config).
; Has full service rights but is NOT root.
; Port 9000 is internal (webpanel_net only), never exposed to the customer.
; www-data (customer SSH/PHP) has zero access to this process or its files.
command=python3 /container_api.py
user=gnukontrolr-admin
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/container-api.err
stdout_logfile=/var/log/supervisor/container-api.out
priority=15
